# root/etc/services.d/stubby/run
#!/usr/bin/execlineb -P

s6-setuidgid stubby

/bin/stubby -C /etc/stubby/stubby.yml
