# root/etc/services.d/knot/run
#!/usr/bin/execlineb -P

/sbin/knotd -c /etc/knot/knot.conf