# root/etc/services.d/kea-dhcp6/run
#!/usr/bin/execlineb -P

# Since the DHCPv6 server opens privileged ports, it requires root access. This daemon must be run as root.

/sbin/kea-dhcp6 -c /etc/kea/kea-dhcp6.conf
