# root/etc/services.d/kea-dhcp6/run
#!/usr/bin/execlineb -P

# this fails
# s6-setuidgid kea

/sbin/kea-dhcp6 -c /etc/kea/kea-dhcp6.conf
