# root/etc/services.d/kea-dhcp4/run
#!/usr/bin/execlineb -P

s6-setuidgid kea

/bin/kea-dhcp4 -c /etc/kea-dhcp4.conf
