# root/etc/services.d/kea-dhcp-ddns/run
#!/usr/bin/execlineb -P

# this fails
# s6-setuidgid kea

/sbin/kea-dhcp-ddns -c /etc/kea/kea-dhcp-dns.conf
