# root/etc/services.d/kea-dhcp-ddns/run
#!/usr/bin/execlineb -P

# Since the DHCPv4 server opens privileged ports, it requires root access. This daemon must be run as root.

/sbin/kea-dhcp-ddns -c /etc/kea/kea-dhcp-ddns.conf
